# kubernetes/etc/profile.d/kube.sh
PATH=/k8s/kubernetes/bin:$PATH:$HOME/bin
