# others/02/kafka-flume-hive/kafka-flume-hive/file/etc/profile.d/hive.sh
export HADOOP_HOME=/opt/hadoop/hadoop-2.9.2
export HADOOP_CONF_DIR=/opt/hadoop/hadoop-2.9.2/conf
export HIVE_HOME=/opt/hive/apache-hive-2.3.5-bin
export HIVE_CONF_DIR=/opt/hive/apache-hive-2.3.5-bin/conf
export PATH=$PATH:$HIVE_HOME/bin:$HADOOP_HOME/bin:$HADOOP_HOME/sbin
